# docking.service
[Unit]
Description=Dual monitor docking

[Service]
User=matt
ExecStart=/usr/local/bin/monitor.sh

[Install]
WantedBy=multi-user.target
